# roles/web/templates/youdo.conf
[program:youdo]
command = /home/alex/.virtualenvs/youdo/bin/gunicorn "app:create_app('app.settings')"
directory = /home/alex/www/youdo/server
user = alex
autostart = true
